# conf/crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * echo $(date "+\%Y-\%m-\%d \%H:\%M:\%S") > /proc/1/fd/2
